# master-ai-telegram.service
[Unit]
Description=Master AI Telegram Bot
After=network-online.target master-ai.service
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/master_ai
Environment=PATH=/home/pi/master_ai/venv/bin:/usr/local/bin:/usr/bin:/bin
EnvironmentFile=/home/pi/master_ai/.env
ExecStart=/home/pi/master_ai/venv/bin/python3 /home/pi/master_ai/telegram_bot.py
Restart=always
RestartSec=10
StandardOutput=append:/home/pi/master_ai/telegram_bot.log
StandardError=append:/home/pi/master_ai/telegram_bot.log

[Install]
WantedBy=multi-user.target
